# .dotfiles/install/by-path/etc/udev/rules.d/99-lowbat.rules
SUBSYSTEM=="power_supply", \
ATTR{capacity_level}=="Low", \
RUN+="/etc/udev/rules.d/99-lowbat.sh"
