# .dotfiles/install/by-path/etc/profile.d/cabal.sh
cabalBin="${CABAL_DIR:-.cabal}/bin"

[ -d "$cabalBin" ] && append_path "$cabalBin"

export PATH
